# systemd/.config/systemd/user/auto-hibernate.service
[Unit]
Description=Auto Hibernate Service

[Service]
ExecStart=/bin/bash %h/bin/auto-hibernate

[Install]
WantedBy=default.target
